# firmware/_antMiner_openwrt20131129.bin.extracted/squashfs-root/usr/lib/opkg/info/pic-update.list
/etc/config/miner_pic.hex
/usr/sbin/pic-update
/etc/init.d/pic-update
/etc/config/pic-update
