# firmware/_antMiner_S320140721.bin.extracted/squashfs-root/usr/lib/opkg/info/cgminer.list
/usr/bin/cgminer
/etc/config/cgminer
/usr/bin/cgminer-monitor
/usr/bin/cgminer-api
/etc/config/asic-freq
/etc/init.d/cgminer
